# config/i-sibal-stable.service
[Unit]
Description=i-Sibal AirPlay Receiver (Stable Mode)
After=network-online.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service

[Service]
Type=simple
Environment="GST_PLUGIN_FEATURE_RANK=v4l2slh264dec:MAX"
# 안정성 모드: sync=true로 오디오/비디오 타이밍을 엄격히 맞추고, max-lateness 50ms로 일시적 간섭 시 프레임 드랍을 막습니다.
ExecStart=/usr/local/bin/uxplay -n "i-Sibal" -vs "kmssink sync=true max-lateness=50000000" -as alsasink
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
